# deploy/systemd/agent-services.slice
[Unit]
Description=Background agent services

[Slice]
CPUAccounting=true
MemoryAccounting=true
IOAccounting=true
TasksAccounting=true
CPUWeight=25
IOWeight=25
